# supervisord.conf
[supervisord]
nodaemon=true

[program:unifi]
command=/usr/lib/unifi/bin/unifi.init start
