# supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod --dbpath /data --fork

[program:unifi]
command=/usr/lib/unifi/bin/unifi.init start
